# univention-node-exporter/debian/univention-node-exporter.cron.d
PATH=/sbin:/bin:/usr/sbin:/usr/bin
*/20 * * * * root /usr/lib/univention-node-exporter/univention-metrics-server.py

